# tarot-bot.service
[Unit]
Description=Tarot Bot Telegram Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=artlvr
Group=artlvr
WorkingDirectory=/home/artlvr/tarot-bot
Environment=PYTHONPATH=/home/artlvr/tarot-bot
Environment=VIRTUAL_ENV=/home/artlvr/tarot-bot/venv
Environment=PATH=/home/artlvr/tarot-bot/venv/bin:$PATH
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/artlvr/tarot-bot/.env

# Запуск с явным указанием интерпретатора и полными путями
ExecStart=/home/artlvr/tarot-bot/venv/bin/python3 /home/artlvr/tarot-bot/manage.py start

# Настройки перезапуска
Restart=always
RestartSec=5
StartLimitBurst=0

# Настройки процесса
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed

# Настройки для предотвращения зомби-процессов
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
