# tarot-bot.service
[Unit]
Description=Tarot Bot Telegram Service
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=artlvr
Group=artlvr
WorkingDirectory=/home/artlvr/tarot-bot

# Environment setup
Environment=PYTHONPATH=/home/artlvr/tarot-bot
Environment=VIRTUAL_ENV=/home/artlvr/tarot-bot/venv
Environment=PATH=/home/artlvr/tarot-bot/venv/bin:$PATH
Environment=PYTHONUNBUFFERED=1

# Load environment variables from file
EnvironmentFile=/home/artlvr/tarot-bot/.env

# Start command
ExecStart=/home/artlvr/tarot-bot/venv/bin/python3 /home/artlvr/tarot-bot/manage.py start

# Restart configuration
Restart=always
RestartSec=5
StartLimitBurst=0

# Process management
TimeoutStartSec=30
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
