# tarot-bot.service
[Unit]
Description=Tarot Bot Telegram Service
After=network.target

[Service]
Type=simple
User=artlvr
Group=artlvr
WorkingDirectory=/home/artlvr/tarot-bot
Environment=PYTHONPATH=/home/artlvr/tarot-bot
Environment=VIRTUAL_ENV=/home/artlvr/tarot-bot/venv
Environment=PATH=/home/artlvr/tarot-bot/venv/bin:$PATH
EnvironmentFile=/home/artlvr/tarot-bot/.env
ExecStart=/home/artlvr/tarot-bot/venv/bin/python3 manage.py start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
